# launchd/ai.agento-manager.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Comment</key>
	<string>Consolidated Agento Manager - manages all agento processes.
Replaces 3 separate plists:
- ai.agento.dashboard.plist
- ai.agento.backfill.plist
- ai.agento.notifier.plist</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>@HOME@</string>
		<key>PATH</key>
		<string>@HOME@/.nvm/versions/node/current/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:@HOME@/bin:/opt/homebrew/bin</string>
		<key>PYTHONUNBUFFERED</key>
		<string>1</string>
		<key>AO_CONFIG_PATH</key>
		<string>@HOME@/agent-orchestrator.yaml</string>
		<key>AO_DASHBOARD_PROJECTS</key>
		<string></string>
		<key>OPENCLAW_AGENTO_DASHBOARD_PROJECTS</key>
		<string></string>
		<!-- GITHUB_TOKEN, GH_TOKEN, OPENCLAW_SLACK_BOT_TOKEN omitted — set in openclaw.json or shell environment -->
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>ai.agento.manager</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-c</string>
		<string>source @HOME@/.bashrc

# Dashboard - starts all project agents configured via environment
echo "[agento-manager] Starting dashboard..."
DASHBOARD_PROJECTS="${AO_DASHBOARD_PROJECTS:-${OPENCLAW_AGENTO_DASHBOARD_PROJECTS:-}}"
if [[ -z "$DASHBOARD_PROJECTS" ]]; then
    echo "[agento-manager] No dashboard projects configured. Set AO_DASHBOARD_PROJECTS or OPENCLAW_AGENTO_DASHBOARD_PROJECTS."
else
    IFS=',' read -r -a DASHBOARD_PROJECT_LIST &lt;&lt;&lt; "$DASHBOARD_PROJECTS"
    for p in "${DASHBOARD_PROJECT_LIST[@]}"; do
        [[ -n "$p" ]] || continue
        @HOME@/bin/ao start "$p" --no-dashboard 2>/dev/null || true
    done
fi

# Backfill - periodic sync
echo "[agento-manager] Starting backfill worker..."
while true; do
    @HOME@/.openclaw/scripts/ao-backfill.sh
    sleep 900  # 15 minutes
done &amp;
BACKFILL_PID=$!

# Kill backfill loop when this script exits (covers notifier crash or clean exit)
# First kill any child processes (ao-backfill.sh, sleep), then the subshell
trap 'pkill -P $BACKFILL_PID 2>/dev/null; kill $BACKFILL_PID 2>/dev/null; wait $BACKFILL_PID 2>/dev/null' EXIT

# Notifier - event-driven notifications
echo "[agento-manager] Starting notifier..."

# Notifier - event-driven notifications (foreground; script exits when it exits)
# Path to notifier script - adjust if your setup differs
NOTIFIER_SCRIPT="${AO_NOTIFIER_SCRIPT:-$HOME/.openclaw/scripts/agento-notifier.py}"
if [[ -f "$NOTIFIER_SCRIPT" ]]; then
    echo "[agento-manager] Starting notifier from $NOTIFIER_SCRIPT..."
    /opt/homebrew/bin/python3 "$NOTIFIER_SCRIPT"
else
    echo "[agento-manager] ERROR: Notifier script not found at $NOTIFIER_SCRIPT"
    exit 1
fi</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/tmp/agento-manager.err.log</string>
	<key>StandardOutPath</key>
	<string>/tmp/agento-manager.log</string>
	<key>WorkingDirectory</key>
	<string>@HOME@</string>
</dict>
</plist>
