# launchd/com.jleechanorg.spend-alert-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jleechanorg.spend-alert-daily</string>

    <key>Comment</key>
    <string>Repo: https://github.com/jleechanorg/smartclaw — Daily GH Actions / optional Gemini &amp; GCP spend (gh api billing). Install: cp scripts/spend-alert-daily.sh ~/.smartclaw/scripts/; cp launchd/com.jleechanorg.spend-alert-daily.plist ~/Library/LaunchAgents/; launchctl bootstrap gui/$UID ~/Library/LaunchAgents/com.jleechanorg.spend-alert-daily.plist</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>${HOME}/.smartclaw/scripts/spend-alert-daily.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>.smartclaw/logs/spend-alert-daily.log</string>

    <key>StandardErrorPath</key>
    <string>.smartclaw/logs/spend-alert-daily.error.log</string>

    <key>WorkingDirectory</key>
    <string>.smartclaw</string>

    <key>ProcessType</key>
    <string>Background</string>

    <key>ThrottleInterval</key>
    <integer>86400</integer>
</dict>
</plist>
