# launchd/com.jleechan.ai-reviewer-stress-test.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jleechan.ai-reviewer-stress-test</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd ~/projects_reference/smartclaw && scripts/ai-reviewer-stress-test.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>14400</integer>
    <!-- 14400 seconds = 4 hours -->

    <key>StandardOutPath</key>
    <string>/tmp/ai-reviewer-stress-test.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/ai-reviewer-stress-test.err.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/jleechan</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>
