# launchd/smartclaw.schedule.orch-health-weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.smartclaw.schedule.orch-health-weekly</string>
    <key>ProgramArguments</key>
    <array>
        <string>@PYTHON3_PATH@</string>
        <string>-m</string>
        <string>orchestration.cron_runner</string>
        <string>--notify</string>
    </array>
    <key>WorkingDirectory</key>
    <string>@HOME@/.smartclaw</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>@HOME@/.smartclaw/src</string>
        <key>HOME</key>
        <string>@HOME@</string>
        <key>PATH</key>
        <string>@HOME@/.cargo/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/openclaw/orch-health-weekly.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/openclaw/orch-health-weekly.err</string>
</dict>
</plist>
