# launchd/smartclaw.gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Comment</key>
	<string>OpenClaw Gateway (tokens live in openclaw.json, not plist env)</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>GOOGLE_CLOUD_PROJECT</key>
		<string>infinite-zephyr-487405-d0</string>
		<key>HOME</key>
		<string>/Users/jleechan</string>
		<key>MISSION_CONTROL_BASE_URL</key>
		<string>http://localhost:9010</string>
		<key>MISSION_CONTROL_BOARD_ID</key>
		<string>aa68f729-d5e0-4d44-8c99-51bcebc0b8bc</string>
		<key>NODE_EXTRA_CA_CERTS</key>
		<string>/etc/ssl/cert.pem</string>
		<key>NODE_USE_SYSTEM_CA</key>
		<string>1</string>
		<key>OPENCLAW_GATEWAY_PORT</key>
		<string>18789</string>
		<key>OPENCLAW_LAUNCHD_LABEL</key>
		<string>ai.smartclaw.gateway</string>
		<key>OPENCLAW_RAW_STREAM</key>
		<string>1</string>
		<key>OPENCLAW_RAW_STREAM_PATH</key>
		<string>/tmp/openclaw/raw-stream.jsonl</string>
		<key>OPENCLAW_SERVICE_KIND</key>
		<string>gateway</string>
		<key>OPENCLAW_SERVICE_MARKER</key>
		<string>openclaw</string>
		<key>OPENCLAW_SERVICE_VERSION</key>
		<string>2026.3.8</string>
		<key>OPENCLAW_SYSTEMD_UNIT</key>
		<string>openclaw-gateway.service</string>
		<key>OPENCLAW_WINDOWS_TASK_NAME</key>
		<string>OpenClaw Gateway</string>
		<key>PATH</key>
		<string>${HOME}/.local/bin:${HOME}/.npm-global/bin:${HOME}/bin:${HOME}/.volta/bin:${HOME}/.asdf/shims:${HOME}/.bun/bin:${HOME}/Library/Application Support/fnm/aliases/default/bin:${HOME}/.fnm/aliases/default/bin:${HOME}/Library/pnpm:${HOME}/.local/share/pnpm:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
		<key>SHELL</key>
		<string>/opt/homebrew/bin/bash</string>
		<key>TMPDIR</key>
		<string>/var/folders/j0/byd1z6px50v88lf679bgt0h00000gn/T/</string>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>ai.smartclaw.gateway</string>
	<key>ProgramArguments</key>
	<array>
		<string>@HOME@/.nvm/versions/node/v22.22.0/bin/node</string>
		<string>@HOME@/.nvm/versions/node/v22.22.0/lib/node_modules/openclaw/dist/index.js</string>
		<string>gateway</string>
		<string>--port</string>
		<string>18789</string>
		<string>--allow-unconfigured</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>${HOME}/.smartclaw/logs/gateway.err.log</string>
	<key>StandardOutPath</key>
	<string>${HOME}/.smartclaw/logs/gateway.log</string>
	<key>ThrottleInterval</key>
	<integer>30</integer>
	<key>Umask</key>
	<integer>63</integer>
</dict>
</plist>
