# launchd/smartclaw.antig-cmux-loop-kill.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.smartclaw.antig-cmux-loop-kill</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>launchctl remove ai.smartclaw.antig-cmux-loop</string>
    </array>
    <key>StartInterval</key>
    <integer>28800</integer>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>${HOME}/.smartclaw/logs/antig-cmux-loop-kill.log</string>
    <key>StandardErrorPath</key>
    <string>${HOME}/.smartclaw/logs/antig-cmux-loop-kill.err</string>
</dict>
</plist>
