# launchd/smartclaw.health-check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.openclaw.health-check</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/jleechan/.openclaw/health-check.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/jleechan</string>
    <key>PATH</key>
    <string>/Users/jleechan/.nvm/versions/node/current/bin:/Users/jleechan/.nvm/versions/node/v22.22.0/bin:/Users/jleechan/Library/pnpm:/Users/jleechan/.bun/bin:/Users/jleechan/.local/bin:/Users/jleechan/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>OPENCLAW_HEALTH_CURL_TIMEOUT</key>
    <string>8</string>
    <key>OPENCLAW_POST_ACTION_WAIT_SECONDS</key>
    <string>4</string>
    <key>OPENCLAW_ESCALATION_COOLDOWN_SECONDS</key>
    <string>3600</string>
  </dict>

  <key>StartInterval</key>
  <integer>300</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/jleechan/.openclaw/logs/health-check.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/jleechan/.openclaw/logs/health-check.log</string>
</dict>
</plist>
